# showcase/onbox/imzero2-demo.service
[Unit]
Description=imzero2 headless pixel-streaming demo — current release (ADR-0085)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=imzero2
Group=imzero2
Environment=IMZERO2_HEADLESS_LISTEN=127.0.0.1:8089
# A served carrier must survive a viewer clicking a UI close ("Quit"): make an
# app-emitted ViewportCommand::Close a no-op so one viewer can't shut the host
# down for everyone (ADR-0085). The carousel also drops the Quit item entirely
# when IMZERO2_HEADLESS_LISTEN is set; this is the defence-in-depth layer.
Environment=IMZERO2_HEADLESS_IGNORE_CLOSE=1
Environment=IMZERO2_HEADLESS_FPS=30
Environment="IMZERO2_HEADLESS_ENCODER_ARGS=-c:v libx264 -preset veryfast -tune zerolatency -bf 0 -g 100000"
Environment=LIBGL_ALWAYS_SOFTWARE=1
ExecStart=/opt/imzero2/run-current.sh
# always (not on-failure): the host can exit cleanly (status 0) on a UI-driven
# viewport-close, which on-failure would NOT restart — that took the demo down
# for 10h once (ADR-0085). With IMZERO2_HEADLESS_IGNORE_CLOSE that exit is now
# rare, so this is the operational safety net. systemd's default start-limit
# still parks a genuine crash-loop (e.g. a broken release) instead of spinning.
Restart=always
RestartSec=2
# Light hardening (the carrier needs no privilege). Tighten further
# (ProtectSystem=strict, etc.) once the FFFI2 shm path is confirmed.
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
